# udev/rules.d/10-dock.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0.0", ENV{XAUTHORITY}="/home/ato/.Xauthority", RUN+="/home/ato/1/dotfiles/dock.sh"

