# udev/rules.d/30-keyboard.rules
KERNEL=="card0", SUBSYSTEM=="drm", RUN+="xmodmap ~/.Xmodmap"

