# udev/rules.d/81-dock.rules
SUBSYSTEM=="usb", ACTION=="add|remove", ENV{ID_VENDOR_ID}=="17ef", ENV{ID_MODEL_ID}=="100f" RUN+="/usr/sbin/dock"
