# Library/LaunchAgents/src.sshfs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>src.sshfs</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/sshfs</string>
		<string>ato@delcorf.sny.no:/src</string>
		<string>/Volumes/src</string>
		<string>-oidmap=user,exec,allow_other,reconnect,auto_cache,default_permissions,volname=src</string>
		<string>-f</string>
	</array>

	<key>RunAtLoad</key>
	<true />

        <key>StandardErrorPath</key>
        <string>/Users/am/Desktop/sshfs.err.log</string>
        <key>StandardOutPath</key>
        <string>/Users/am/Desktop/sshfs.out.log</string>
</dict>
</plist>
